# supervisord.conf
[program:sample_stream]
command=tl sample

[supervisord]
logfile=/var/log/supervisord.log
loglevel=debug
