# supervisord.conf
[program:sample]
command=tl sample

[program:filter_sad]
command=tl filter sad data/sad.json

[program:filter_happy]
command=tl filter happy data/happy.json

[supervisord]
logfile=/var/log/supervisord.log
loglevel=debug

[supervisorctl]
serverurl=unix:///tmp/supervisord.sock

[unix_http_server]
file=/tmp/supervisord.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface
