# supervisord.conf
[program:sample_stream]
command=tl stream

[supervisord]
logfile=/var/log/supervisord.log
loglevel=debug
